# use-cases/03-alma-banking-assistant/voice/alma-banking-voice.service
[Unit]
Description=Alma Banking Voice Assistant (AI Bank)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/alma-banking
Environment=PATH=/home/ec2-user/bis-ai-assistant/venv/bin:/usr/local/bin:/usr/bin
Environment=LD_LIBRARY_PATH=/usr/local/lib
Environment=ALMA_SONIC_REGION=eu-north-1
Environment=ALMA_VOICE=matthew
Environment=DB_REGION=me-south-1
Environment=CLUSTER_ARN=arn:aws:rds:me-south-1:519124228967:cluster:aibank-core-banking
Environment=SECRET_ARN=arn:aws:secretsmanager:me-south-1:519124228967:secret:aibank-core-banking-credentials-DEdCPJ
Environment=DB_NAME=corebanking
Environment=BEDROCK_AGENTCORE_MEMORY_ID=alma_banking_assistant_mem-ijns9pFcc6
Environment=MEMORY_REGION=eu-west-1
ExecStart=/home/ec2-user/bis-ai-assistant/venv/bin/python alma_banking_voice.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
